# fs/etc/services.d/prometheus/run
#!/usr/bin/with-contenv sh

s6-setuidgid nonroot

exec /bin/prometheus --config.file=/config/prometheus.yml --storage.tsdb.path=/prometheus $PROMETHEUS_ARGS